# tmux-save.timer
[Unit]
Description=
After=tmux.service

[Timer]
OnCalendar=*:*:0,30

[Install]
WantedBy=multi-user.target
