# tmux-save.service
[Unit]
Description=

[Service]
Type=simple
ExecStart=/home/adion/bin/tmux-save

[Install]
WantedBy=multi-user.target
